# overlay/etc/systemd/system/systemd-timesyncd.service.d/privatetmp.conf
# if PrivateTmp=yes (the default) it fails with "226/NAMESPACE"
[Service]
PrivateTmp=no
